# oscam/config/supervisord.conf
[supervisord]
nodaemon=true

[program:oscam]
command=/var/bin/oscam -c /etc/oscam/
